# root/etc/services.d/snapserver/run
#!/usr/bin/with-contenv sh

/usr/bin/snapserver -c /etc/snapserver.conf

