# root/etc/services.d/scan/run
#!/usr/bin/with-contenv sh

echo "Scanning"
/usr/bin/mopidyctl --config /etc/mopidy/mopidy.conf local scan
sleep 60

